# fabfile/supervisord.conf
[program:%(GROUP)s_%(APP)s]
command=%(ROOT)s/%(APP)s/bin/gunicorn_start.sh                    ; Command to start app
user=%(APP)s                                                          ; User to run as
stdout_logfile=%(ROOT)s/%(APP)s/logs/gunicorn_supervisor.log   ; Where to write log messages
redirect_stderr=true                                                ; Save stderr in the same log

[program:%(GROUP)s_%(APP)s_celery]
command=%(ROOT)s/%(APP)s/bin/celery_start.sh                    ; Command to start app
user=%(APP)s                                                          ; User to run as
stdout_logfile=%(ROOT)s/%(APP)s/logs/celery_supervisor.log   ; Where to write log messages
redirect_stderr=true                                                ; Save stderr in the same log
